# var/lib/dpkg/info/mpd.list
/.
/usr
/usr/bin
/usr/bin/mpd
/usr/share
/usr/share/doc/mpd
/usr/share/doc/mpd/examples
/usr/share/doc/mpd/examples/mpdconf.example.gz
/usr/share/doc/mpd/examples/mpd.desktop
/usr/share/doc/mpd/changelog.Debian.gz
/usr/share/doc/mpd/changelog.gz
/usr/share/doc/mpd/user-manual.html
/usr/share/doc/mpd/AUTHORS
/usr/share/doc/mpd/README.Debian
/usr/share/doc/mpd/NEWS.Debian.gz
/usr/share/doc/mpd/copyright
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/mpd.conf.5.gz
/usr/share/man/man1
/usr/share/man/man1/mpd.1.gz
/usr/share/doc-base
/usr/share/doc-base/mpd
/usr/lib
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/mpd.conf
/etc
/etc/mpd.conf
/etc/logrotate.d
/etc/logrotate.d/mpd
/etc/default
/etc/default/mpd
/etc/init.d
/etc/init.d/mpd
/var
/var/log
/var/log/mpd
/var/lib
/var/lib/mpd
/var/lib/mpd/playlists
/var/lib/mpd/music
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/mpd.service
