# var/lib/dpkg/info/passwd.conffiles
/etc/default/useradd
/etc/cron.daily/passwd
/etc/pam.d/chpasswd
/etc/pam.d/chfn
/etc/pam.d/newusers
/etc/pam.d/chsh
/etc/pam.d/passwd
