# var/lib/dpkg/info/sudo.conffiles
/etc/sudoers.d/README
/etc/init.d/sudo
/etc/sudoers
/etc/pam.d/sudo
