# var/lib/dpkg/info/winbind.list
/etc
/etc/default
/etc/default/winbind
/etc/init.d
/etc/init.d/winbind
/etc/logrotate.d
/etc/logrotate.d/winbind
